# conn-exporter.service
[Unit]
Description=Network Connection Exporter
Documentation=https://github.com/bjin01/conn-exporter
After=network.target
Wants=network.target

[Service]
Type=simple
ExecStart=/usr/bin/conn-exporter
Restart=always
RestartSec=5

# Run as root to access network information
User=root
Group=root

# Environment
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Working directory
WorkingDirectory=/

[Install]
WantedBy=multi-user.target